# langate2000-supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true

[program:gunicorn]
directory=/app/langate/
command=gunicorn langate.wsgi
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:langate2000-netcontrol]
directory=/app/langate2000-netcontrol/
command=python3 netcontrol.py
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
